# share/autostart/io.github.hyperx_saga_control.desktop.in
[Desktop Entry]
Type=Application
Name=HyperX Saga Control
Comment=Start HyperX Saga Control in the system tray
Exec=__BIN__ --start-minimized
Icon=hyperx-saga-control
Terminal=false
X-KDE-autostart-after=panel
X-GNOME-Autostart-enabled=true
